# backup_service.service
[Unit]
Description=FastAPI Backup + YubiKey Warm-Up Service
After=network-online.target
Wants=network-online.target

[Service]
User=fraser
WorkingDirectory=/var/home/fraser/backup_service

# Start smartcard and GPG infrastructure
ExecStartPre=/usr/bin/gpgconf --launch gpg-agent

Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus


# Make uv visible to systemd
Environment="PATH=/var/home/fraser/.cargo/bin:/usr/local/bin:/usr/bin:/bin"

# Use bash -lc to load your user environment
ExecStart=/usr/bin/bash -lc "/var/home/fraser/.cargo/bin/uv run uvicorn main:app --host 127.0.0.1 --port 8000"

Restart=on-failure
RestartSec=5
ProtectHome=no

[Install]
WantedBy=multi-user.target
